# .config/LSD/test/99-monitor-mode.rules
ACTION=="add", SUBSYSTEM=="net", ATTR{idVendor}=="2357", ATTR{idProduct}=="010c", RUN+=""

